# gui-src/gl-sdk4-tailscale/data/etc/hotplug.d/net/19-tailscale-net
#!/bin/sh

if [ "$INTERFACE" = "tailscale0" ] && [ "$ACTION" = "add" ];then
    /usr/bin/gl_tailscale  set_route &
fi
